# donkey-kong.service
[Unit]
Description=Donkey Kong SSH Core Daemon
After=network.target

[Service]
ExecStart=/usr/bin/npx tsx /opt/donkey-kong/daemon.ts
WorkingDirectory=/opt/donkey-kong
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
Restart=always
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
